# cryptoconvbot.service
# Put this in /lib/systemd/system/cryptoconvbot.service
# Set chmod 644 /lib/systemd/system/cryptoconvbot.service

[Unit]
Description=CryptoConvBot
After=syslog.target
After=network.target

[Service]
# Start main service
ExecStart=/usr/bin/python /bots/CryptoConvBot/CryptoConvBot.py > /bots/CryptoConvBot/service.log 2>&1
Restart=always
RestartSec=600

[Install]
WantedBy=multi-user.target
Alias=CryptoConv.bot
